# display-manager.service
[Unit]
Description=Display Manager Container
After=docker.service network-online.target
Requires=docker.service
Before=kiosk.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/jkumar/github/display-manager
ExecStart=/usr/bin/docker compose -f /home/jkumar/github/display-manager/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /home/jkumar/github/display-manager/docker-compose.prod.yml down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
