# numpad-watcher/numpad-watcher.service
[Unit]
Description=Numpad Watcher for Display Manager Kiosk
After=display-manager.service
Requires=display-manager.service

[Service]
Type=simple
User=jkumar
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/jkumar/github/display-manager/numpad-watcher/watcher.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
